# Main/Logged/etc/crontabs/root
*/2 * * * * /bin/sh /root/test.sh
49 23 * * * /bin/sh /root/stat.sh
